# recipes-multimedia/janus-gateway/files/janus-gateway.service
[Unit]
Description=Janus WebRTC Gateway
After=network.target
Wants=network.target

[Service]
Type=simple
SyslogIdentifier=janus
ExecStartPre=/bin/bash -c "JANUS_ADMIN_SECRET=$(cat ${CREDENTIAL_DIRECTORY}/janus-edge-admin-secret) \
    JANUS_SERVER_NAME=$(cat /etc/hostname) \
    janus-envsubst-on-templates"
ExecStart=/usr/bin/janus
Restart=always
RestartSec=60
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
