# meta-printnanny/recipes-core/printnanny-core-apps/files/printnanny-online.target
[Unit]
Description=PrintNanny OS is online & ready to use
After=multi-user.target
