# meta-printnanny/recipes-core/printnanny-core-apps/files/printnanny-shutdown.service
[Unit]
Description=PrintNanny shutdown service hook

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/true
ExecStart=/usr/bin/printnanny -v os shutdown

[Install]
WantedBy=multi-user.target
