# recipes-core/bitsy-tweaks/files/getty@tty1.service.d/50-noclear.conf
[Unit]
After=cloud-init-local.service
[Service]
TTYVTDisallocate=no
