# meta-printnanny/recipes-core/printnanny-vision/files/printnanny-vision.service
[Unit]
Description=PrintNanny Cam Service
After=network-online.target printnanny-cloud-sync.service
Wants=network-online.target printnanny-cloud-sync.service 
StartLimitInterval=60
StartLimitBurst=3

[Service]
Type=simple
SyslogIdentifier=printnanny-vision
ExecStart=/usr/bin/printnanny-vision --config /var/run/printnanny/printnanny-vision.toml
Restart=on-failure
RestartSec=10
User=printnanny

[Install]
WantedBy=multi-user.target
