# meta-printnanny/recipes-core/printnanny-vision/files/printnanny-vision.target
[Unit]
Description=Systemd unit to group all PrintNanny Vision services
Wants=gstd.service printnanny-vision.service

[Install]
WantedBy=multi-user.target