# recipes-core/bitsy-tweaks/files/bitsy-growfs.service
[Unit]
Description=Resize last partition/fs to fill disk
Before=local-fs-pre.target
DefaultDependencies=no

[Service]
Type=oneshot
TimeoutSec=infinity
ExecStart=/usr/sbin/bitsy-growfs
ExecStart=/bin/systemctl --no-reload disable %n
StandardOutput=journal+console

[Install]
RequiredBy=local-fs-pre.target
