# recipes-core/bitsy-tweaks/files/bitsy-overlayfs-all.target.in
[Unit]
Description=Groups all overlays required by {PN} in one unit
After={ALL_OVERLAYFS_UNITS}
Requires={ALL_OVERLAYFS_UNITS}

[Install]
WantedBy=local-fs.target
