# recipes-core/bitsy-tweaks/files/bitsy-growfs.service
[Unit]
Description=Resize last partition/fs to fill disk
DefaultDependencies=no
Conflicts=shutdown.target
After=systemd-remount-fs.service
ConditionPathExists=!/etc/bitsy/growfs.done

[Service]
Type=oneshot
RemainAfterExit=yes
StandardOutput=journal+console
ExecStart=/bin/bash -c "yes | parted ---pretend-input-tty /dev/mmcblk0 resizepart 3 100%"
ExecStart=/sbin/resize2fs /dev/mmcblk0p3
ExecStart=/bin/bash -c "echo $(date +%s) > /etc/bitsy/growfs.done"
