# meta-printnanny/recipes-extended/cloud-init/cloud-init/cloud-init.target
[Unit]
Description=Cloud-init target
Before=multi-user.target

[Install]
WantedBy=multi-user.target
